# deploy/systemd/xray-control.service
[Unit]
Description=Xray Control API
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/home/stablespace/projects/backend
EnvironmentFile=/home/stablespace/projects/backend/.env
ExecStart=/usr/bin/env uvicorn app.main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=3
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
